# net-misc/forked-easyfind-client/files/easyfind-client.initd
#!/sbin/openrc-run

NAME=easyfind-client
DAEMON=/opt/bubba/sbin/efd
PIDFILE=/run/bubba-easyfind.pid

start() {
        ebegin "Starting ${NAME}"
		start-stop-daemon --start --quiet --background --exec ${DAEMON}
        eend $?
}

stop() {
        ebegin "Stopping ${NAME}"
	        start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
        eend $?
}
